# rstudio/s6/s6-overlay/s6-rc.d/rstudio/finish
#!/command/execlineb -P

rstudio-server stop
