# codeserver/s6/s6-overlay/s6-rc.d/code-server/run
#!/command/execlineb -P

with-contenv
/usr/bin/code-server \
  --bind-addr 0.0.0.0:8888 \
  --disable-telemetry \
  --auth none
