# rstudio/s6/s6-overlay/s6-rc.d/rstudio-env-fix/up
with-contenv

/etc/s6-overlay/scripts/rstudio-env-fix
